# systemd/ydotoold.service
[Unit]
Description=ydotool virtual-input daemon (for ntok text injection)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ydotoold --socket-path=%t/.ydotool_socket --socket-perm=0600
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
